# deploy/etc/systemd/system/bluetooth.service.d/gamepadserver.conf
# GamePadServer: bluetoothd overrides for HID emulation
#
# --compat       Enable legacy SDP socket interface (needed for RegisterProfile)
# --noplugin=input  Prevent BlueZ from intercepting HID L2CAP connections
#
# Managed by: deploy/setup-host.sh
# See also: .claude/docs/bluetooth/pitfalls.md

[Service]
ExecStart=
ExecStart=/usr/libexec/bluetooth/bluetoothd --compat --noplugin=input
